# ems-api.service.template
[Unit]
Description=EMS REST API Server
Documentation=file://__EMS_DIR__/api_server.py
After=network-online.target postgresql.service ems-air-quality.service ems-seismic.service
Wants=network-online.target
PartOf=ems.target

[Service]
Type=simple
User=__EMS_USER__
WorkingDirectory=__EMS_DIR__
ExecStart=/usr/bin/python3 __EMS_DIR__/api_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ems-api

[Install]
WantedBy=multi-user.target ems.target
